# crontab
*/2 * * * * echo "Hello Cron!" > /dev/stdout
